# default/ncr-bootstrap.service
[Unit]
Description=NCR Bootstrap Service
After=network.Target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
